# Archive/void-files/etc/skel/.config/autostart/X_PostSettings.desktop
[Desktop Entry]
Encoding=UTF-8
Type=Application
Name=X_PostSettings
Exec=/usr/local/bin/void-scripts/xpostsettings.sh
OnlyShowIn=XFCE;
StartupNotify=false
Terminal=false
Hidden=false